# syncthing/supervisord.conf
[supervisord]
nodaemon=true

[program:syncthing]
command=/syncthing -home=/config
environment=STNORESTART=yes
autorestart=true
redirect_stderr=true
user=syncthing
